# sv/kafka/run
#!/bin/bash

#sv start socklog-unix || exit 1

source /etc/envvars
env

cd /kafka/bin

exec 2>&1
exec ./kafka-server-start.sh ../config/server.properties 
